# mkosi.profiles/rootfs-embedded/mkosi.extra/usr/lib/systemd/system/run-rootfs.mount
[Unit]
After=modprobe@erofs.service
Requires=modprobe@erofs.service

[Mount]
What=/rootfs.ero
Where=/run/rootfs
Type=erofs
