# mkosi.extra/usr/lib/systemd/system/buffyboard.service.d/start-after-getty.conf
[Unit]
Requires=getty@tty1.service
After=getty@tty1.service
