# mkosi.extra/usr/lib/systemd/system/buffyboard.service.d/wait-for-login.conf
[Unit]
Requires=user@0.service
After=user@0.service
