# mkosi.extra/usr/lib/systemd/system/buffyboard.service.d/wait.conf
# If buffyboard starts up too fast it doesn't properly resize the terminal
[Service]
ExecStartPre=/usr/bin/sleep 2
