# mkosi.profiles/buffyboard/mkosi.extra/etc/profile.d/buffyboard.sh
if [[ "$XDG_SEAT" == *"seat"* ]]; then
    systemctl start buffyboard
fi
