# mkosi.extra/usr/lib/systemd/system/systemd-vconsole-setup.service.d/hidpi-consolefont.conf
[Service]
ExecStartPre=/usr/bin/setfont -d
